# sabnzbd/rootfs/etc/s6-overlay/s6-rc.d/sabnzbd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: SABnzbd
# Runs SABnzbd
# ==============================================================================

bashio::log.info 'Starting SABnzbd...'

exec python3 /opt/SABnzbd.py \
    --config-file /config \
    --server 0.0.0.0 \
    --console \
    --disable-file-log
